# conf/supervisor/worker.conf
[program:long_script]
command=/usr/bin/php /var/www/current/bin/worker.php
autostart=true
autorestart=true
stderr_logfile=/var/log/supervisor/worker.err.log
stdout_logfile=/var/log/supervisor/worker.out.log
